# crontab/apt-mirror
#
# Regular cron jobs for the apt-mirror package
#
# 0 4 * * * - Run command at every 04:00

0 */6 * * *   /usr/bin/apt-mirror > /var/spool/apt-mirror/var/cron.log



